# ansible/roles/docker/files/etc/profile.d/99-docker.sh
export CONTAINER_RUNTIME_ENDPOINT="unix:///var/run/dockershim.sock"
