# ansible/roles/auditd/files/etc/systemd/system/auditd-rotate.service
[Unit]
Description=Rotate auditd logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/service auditd rotate

